# puppet/yakkl_ops/files/cron.d/log2yakkl
MAILTO=root

* * * * *   yakkl cd /home/yakkl/deployments/current && python3 ./yakkl-current-venv/share/yakkl/integrations/log2yakkl/log2yakkl
