# puppet/yakkl/files/cron.d/check-analytics-state
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=yakkl

30 * * * * yakkl /home/yakkl/deployments/current/manage.py check_analytics_state
