# puppet/yakkl/files/cron.d/email-mirror
MAILTO=root

* * * * *   yakkl cd /home/yakkl/deployments/current && ./manage.py email_mirror
